# installers/docker/full/image_files/netcentral-transceiver-aprsis
#!/sbin/openrc-run
name="NetCentral APRS-IS Transceiver"
description="NetCentral connectivity to APRS-IS"
retry="60"
pidfile="/run/netcentral-transceiver-aprsis/$RC_SVCNAME.pid"
command="/opt/java/openjdk/bin/java"
command_args="-jar /netcentral/src/transceiver-aprsis.jar </dev/null &>/dev/null &"

depend() {
    need netcentral-server
}
