# installers/docker/split/image_files/netcentral-server
#!/sbin/openrc-run
name="NetCentral Server"
description="NetCentral APRS Nets and more!"
retry="60"
pidfile="/run/netcentral-server/$RC_SVCNAME.pid"
command="/opt/java/openjdk/bin/java"
command_args="-jar /netcentral/bin/netcentral-server.jar </dev/null &>/dev/null &"

depend() {
    need net
}
