# installers/docker/split/image_files/netcentral-ui
#!/sbin/openrc-run
name="NetCentral User Interface"
description="NetCentral UI - APRS Nets and more!"
retry="60"
pidfile="/run/netcentral-ui/$RC_SVCNAME.pid"
# command="cd /netcentral/bin/ui;npm run dev </dev/null &>/dev/null &"

workingdir="/netcentral/bin/ui"

# The command to execute (uses 'npm run dev')
command="/usr/bin/npm"
command_args="run dev"

depend() {
    need net
}
