# installers/docker/split/image_files/osm-proxy-cache
#!/sbin/openrc-run
name="NetCentral OSM Proxy Cache"
description="NetCentral OSM map tile proxy and cache"
retry="60"
pidfile="/run/osm-proxy-cache/$RC_SVCNAME.pid"
command="/opt/java/openjdk/bin/java"
command_args="-jar /netcentral/bin/osm-proxy-cache.jar </dev/null &>/dev/null &"

depend() {
    need net
}
